# services/ddm.service.in
[Unit]
Description=Simple Desktop Display Manager
Documentation=man:ddm(1) man:ddm.conf(5)
Conflicts=getty@tty${DDM_INITIAL_VT}.service
After=systemd-user-sessions.service getty@tty${DDM_INITIAL_VT}.service plymouth-quit.service systemd-logind.service
PartOf=graphical.target
StartLimitIntervalSec=30
StartLimitBurst=2

Requires=seatd.service
Before=seatd.service

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/ddm
Restart=always

CapabilityBoundingSet=CAP_CHOWN CAP_FOWNER CAP_KILL CAP_SETGID CAP_SETUID CAP_SYS_ADMIN CAP_SYS_BOOT CAP_SYS_CHROOT CAP_SYS_TTY_CONFIG
OOMScoreAdjust=-300
Nice=-15
PrivateIPC=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictSUIDSGID=true

[Install]
Alias=display-manager.service
